# packaging/systemd/graywolf.service
[Unit]
Description=Graywolf APRS radio transceiver and gateway
Documentation=https://github.com/chrissnell/graywolf
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/graywolf -config /var/lib/graywolf/graywolf.db -history-db /var/lib/graywolf/graywolf-history.db -tile-cache-dir /var/lib/graywolf/tiles -modem /usr/bin/graywolf-modem -http 0.0.0.0:8080
Restart=on-failure
RestartSec=5

User=graywolf
Group=graywolf
SupplementaryGroups=audio dialout plugdev gpio

# Writable state for SQLite database
StateDirectory=graywolf

# Allow binding to privileged ports
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE

# Filesystem hardening
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/graywolf

# Network
PrivateNetwork=no
RestrictAddressFamilies=AF_INET AF_INET6 AF_UNIX AF_NETLINK

# Kernel hardening
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectKernelLogs=yes
ProtectControlGroups=yes
ProtectClock=yes
ProtectHostname=yes

# Execution hardening
NoNewPrivileges=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
RemoveIPC=yes
SystemCallArchitectures=native

[Install]
WantedBy=multi-user.target
